# deploy/supervisor.conf
[program:task-producer]
command=/path/to/task_queue/venv/bin/python producer.py
directory=/path/to/task_queue
autostart=true
autorestart=true
startsecs=5
stopwaitsecs=10
stdout_logfile=/path/to/task_queue/logs/producer.log
stderr_logfile=/path/to/task_queue/logs/producer.err
environment=PYTHONUNBUFFERED=1

[program:task-monitor]
command=/path/to/task_queue/venv/bin/python monitor.py
directory=/path/to/task_queue
autostart=true
autorestart=true
startsecs=3
stdout_logfile=/path/to/task_queue/logs/monitor.log
stderr_logfile=/path/to/task_queue/logs/monitor.err

[program:task-crawler-ig-full]
command=/path/to/task_queue/venv/bin/python ig_crawler.py --mode full
directory=/path/to/task_queue
autostart=false
autorestart=true
startsecs=15
stopwaitsecs=30
stdout_logfile=/path/to/task_queue/logs/crawler-ig-full.log
stderr_logfile=/path/to/task_queue/logs/crawler-ig-full.err
environment=PYTHONUNBUFFERED=1

[program:task-crawler-ig-incr]
command=/path/to/task_queue/venv/bin/python ig_crawler.py --mode incr
directory=/path/to/task_queue
autostart=false
autorestart=true
startsecs=15
stopwaitsecs=30
stdout_logfile=/path/to/task_queue/logs/crawler-ig-incr.log
stderr_logfile=/path/to/task_queue/logs/crawler-ig-incr.err
environment=PYTHONUNBUFFERED=1

[program:task-crawler-x-full]
command=/path/to/task_queue/venv/bin/python x_crawler.py --mode full
directory=/path/to/task_queue
autostart=false
autorestart=true
startsecs=15
stopwaitsecs=30
stdout_logfile=/path/to/task_queue/logs/crawler-x-full.log
stderr_logfile=/path/to/task_queue/logs/crawler-x-full.err
environment=PYTHONUNBUFFERED=1

[program:task-crawler-x-incr]
command=/path/to/task_queue/venv/bin/python x_crawler.py --mode incr
directory=/path/to/task_queue
autostart=false
autorestart=true
startsecs=15
stopwaitsecs=30
stdout_logfile=/path/to/task_queue/logs/crawler-x-incr.log
stderr_logfile=/path/to/task_queue/logs/crawler-x-incr.err
environment=PYTHONUNBUFFERED=1

[program:task-downloader-ig]
command=/path/to/task_queue/venv/bin/python sub_task_worker.py --mode ig
directory=/path/to/task_queue
autostart=false
autorestart=true
startsecs=5
stdout_logfile=/path/to/task_queue/logs/downloader-ig.log
stderr_logfile=/path/to/task_queue/logs/downloader-ig.err
environment=MAX_TASKS_PER_WORKER=0,PYTHONUNBUFFERED=1

[program:task-downloader-x]
command=/path/to/task_queue/venv/bin/python sub_task_worker.py --mode x
directory=/path/to/task_queue
autostart=false
autorestart=true
startsecs=5
stdout_logfile=/path/to/task_queue/logs/downloader-x.log
stderr_logfile=/path/to/task_queue/logs/downloader-x.err
environment=MAX_TASKS_PER_WORKER=0,PYTHONUNBUFFERED=1
